# deploy/mattube-server.service
[Unit]
Description=mattube server (Drive poller + job processor)
After=network.target

[Service]
Type=simple
User=mattube
WorkingDirectory=/opt/mattube
ExecStart=/usr/bin/mattube-server -c /etc/mattube/config.json
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
